# dot_config/systemd/user/low-battery.service
[Unit]
Description=Low battery warning or suspension
After=systemd-user-sessions.service

[Service]
ExecStart=/home/mberrio/bin/battery
