# config/systemd/user/low-battery.service
[Unit]
Description=Low battery warning or suspension

[Service]
ExecStart=/home/mberrio/new-dotfiles/scripts/battery.sh
